# samba-podman-generated.service
# container-b423f39077a5b38c1710e01d86bf7fc5dad6d21c522b104631978f86a319fc11.service
# autogenerated by Podman 2.0.2
# Wed Jul 15 08:19:26 CDT 2020

[Unit]
Description=Podman container-b423f39077a5b38c1710e01d86bf7fc5dad6d21c522b104631978f86a319fc11.service
Documentation=man:podman-generate-systemd(1)
Wants=network.target
After=network-online.target

[Service]
Environment=PODMAN_SYSTEMD_UNIT=%n
Restart=on-failure
ExecStartPre=/bin/rm -f %t/container-b423f39077a5b38c1710e01d86bf7fc5dad6d21c522b104631978f86a319fc11.pid %t/container-b423f39077a5b38c1710e01d86bf7fc5dad6d21c522b104631978f86a319fc11.ctr-id
ExecStart=/usr/bin/podman run --conmon-pidfile %t/container-b423f39077a5b38c1710e01d86bf7fc5dad6d21c522b104631978f86a319fc11.pid --cidfile %t/container-b423f39077a5b38c1710e01d86bf7fc5dad6d21c522b104631978f86a319fc11.ctr-id --cgroups=no-conmon --replace --privileged --name sammy -d lochnerr/samba-dc
ExecStop=/usr/bin/podman stop --ignore --cidfile %t/container-b423f39077a5b38c1710e01d86bf7fc5dad6d21c522b104631978f86a319fc11.ctr-id -t 10
ExecStopPost=/usr/bin/podman rm --ignore -f --cidfile %t/container-b423f39077a5b38c1710e01d86bf7fc5dad6d21c522b104631978f86a319fc11.ctr-id
PIDFile=%t/container-b423f39077a5b38c1710e01d86bf7fc5dad6d21c522b104631978f86a319fc11.pid
KillMode=none
Type=forking

[Install]
WantedBy=multi-user.target default.target
