# samba-ad-dc.service
[Unit]
Description=Samba Active Directory Domain Controller Service for running samba-dc in container under systemd.
Wants=network-online.target
After=network.target network-online.target

[Service]
Type=simple
ExecStart=/usr/local/bin/samba-run
ExecStop=/bin/echo "shutdown" >/var/lib/samba/cmd
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
